# server/hacktheplanet.service
[Unit]  
Description=dotnetcore 2.1 sample service with generic host  
  
[Service]  

ExecStart=/usr/bin/dotnet HackThePlanet.Host.dll  
WorkingDirectory=/usr/local/sbin/hacktheplanet/
User=root
Group=root
Restart=on-failure  
SyslogIdentifier=hacktheplanet-service  
PrivateTemp=true  
  
[Install]  
WantedBy=multi-user.target